# Cloud_Server/river_flooding_detection.conf
# Script to run the river level detection
# Add this file to /etc/init/

description "Script to run the river level detection"
start on runlevel [2345]
stop on runlevel [06]

respawn

script
    chdir /home/suporte/flash_flooding_alert_system
    . project_env/bin/activate
    python river_level_detection.py
end script
